# infra/ansible/roles/patroni-raft-witness/templates/patroni-raft-witness.service.j2
[Unit]
Description=Patroni raft witness (quorum-only, no PostgreSQL)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{ patroni_venv }}/bin/patroni_raft_controller {{ patroni_raft_witness_config_path }}
KillMode=process
TimeoutSec=30
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
